# docker/supervisor/10-queue.conf
; Laravel default queue worker
[program:queue]
priority=10
command=php /var/www/html/artisan queue:work --tries=1
